# packages/hackinglab-userhome/hackinglab-userhome-0.3/home/hacker/.config/autostart/NotesScreenlet.desktop
[Desktop Entry]
Name=NotesScreenlet
Encoding=UTF-8
Version=1.0
Type=Application
Exec= python -u /usr/share/screenlets/Notes/NotesScreenlet.py
X-GNOME-Autostart-enabled=true
