# packages/hackinglab-autostart/gopenvpn.desktop
[Desktop Entry]
Encoding=UTF-8
Name=gopenvpn
GenericName=OpenVPN GUI
Comment=gopenvpn
Exec=/usr/bin/gopenvpn
Icon=/usr/share/gopenvpn/gopenvpn-icon.png
Terminal=false
Type=Application
